# kong.service
[Unit]
Description=Kong
After=network.target

[Service]
ExecStart=%h/go/bin/kong daemon

[Install]
WantedBy = default.target
